# habarkita.service
# ─────────────────────────────────────────────────────────────
#  Systemd Service — jalankan Flask via Gunicorn
#  Simpan di: /etc/systemd/system/undangankita.service
#
#  Setup:
#    sudo systemctl daemon-reload
#    sudo systemctl enable undangankita
#    sudo systemctl start undangankita
# ─────────────────────────────────────────────────────────────

[Unit]
Description=UndanganKita Flask App
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/undangankita

# Load .env otomatis
EnvironmentFile=/var/www/undangankita/.env

# Gunicorn — lebih production-ready dari flask run
ExecStart=/var/www/undangankita/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --timeout 60 \
    --access-logfile /var/log/undangankita/access.log \
    --error-logfile /var/log/undangankita/error.log \
    app:app

Restart=always
RestartSec=5

# Hardening proses (cegah eskalasi privilege)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/undangankita/data /var/www/undangankita/static/uploads

[Install]
WantedBy=multi-user.target
